# systemd/carbuddy.service
[Unit]
Description=Sentry CarBuddy OBD-II Diagnostic Service
After=obd-connect.service network.target
Requires=obd-connect.service
Wants=network.target

[Service]
Type=simple
ExecStart=/opt/carbuddy/.venv/bin/python /opt/carbuddy/src/main.py
WorkingDirectory=/opt/carbuddy
Restart=on-failure
RestartSec=15
User=root
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/opt/carbuddy/src

[Install]
WantedBy=multi-user.target
